# scripts/com.presence-hub.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.presence-hub.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ilyavereshchagin/Personal/esp32_spt/daemon/run.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ilyavereshchagin/Personal/esp32_spt</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ilyavereshchagin/Library/Logs/presence-hub.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ilyavereshchagin/Library/Logs/presence-hub.err.log</string>
</dict>
</plist>
